# etc/systemd/system/cd-rip.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/run_cd_script.sh
